# submission/week-01/07.service
[Unit]
Description=Hello World oneshot service

[Service]
Type=oneshot
ExecStart=/bin/echo "Hello World"
